# aeternus-os/archiso/airootfs/etc/systemd/system/amnesia-shutdown.service
[Unit]
Description=V0RTEX Amnesia — RAM wipe before poweroff
DefaultDependencies=no
Before=shutdown.target reboot.target halt.target
Requires=shutdown.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/amnesia --confirm
TimeoutStartSec=120
RemainAfterExit=yes

[Install]
WantedBy=halt.target reboot.target shutdown.target
