# aeternus-os/archiso/airootfs/etc/systemd/system/vortex-detect-env.service
[Unit]
Description=V0rtexOS — Detecção de Ambiente (VM vs Hardware Real)
Documentation=file:///usr/local/bin/vortex-detect-env
DefaultDependencies=no
After=systemd-udevd.service local-fs.target
Before=getty@tty1.service display-manager.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/vortex-detect-env
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vortex-detect-env

# Precisa de acesso a /sys, /proc e /run
ProtectSystem=no
PrivateTmp=no

[Install]
WantedBy=sysinit.target
